# nix/packages/openwrt/etc/hotplug.d/iface/70-pbr
#!/bin/sh
#? https://openwrt.org/docs/guide-user/base-system/hotplug
#? https://docs.mossdef.org/pbr/#a-word-about-interface-hotplug-script

[ "$INTERFACE" = "awg0" ] || exit 0
[ "$ACTION" = "ifup" ] || exit 0

logger -t pbr "Reloading $INTERFACE due to $ACTION of $INTERFACE ($DEVICE)"
/etc/init.d/pbr on_interface_reload "$INTERFACE"
