# nix/packages/openwrt/etc/uci-defaults/70-pbr
#!/bin/sh
#? docs.openwrt.melmac.ca/pbr/1.1.8/#Warning:Pleasesetdhcp.lan.force1
uci set dhcp.lan.force='1'
uci commit
