# nix/packages/openwrt/etc/uci-defaults/91-change-root-shell
#!/bin/sh
sed -i 's|/bin/ash|/bin/bash|' /etc/passwd
