# nix/packages/openwrt/etc/uci-defaults/99-custom
#!/bin/sh
uci -q batch << EOI
set system.@system[0].hostname='testap'
commit
EOI
